# chef/cookbooks/zaproxy/templates/default/zaproxy.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides:          zaproxy
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the zaproxy penetration testing tool
# Description:       starts zaproxy using start-stop-daemon
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
NAME=zaproxy
DESC=zaproxy

set -e

. /lib/lsb/init-functions

case "$1" in
  start)
    echo -n "Starting $DESC: "
    start-stop-daemon -b -c zaproxy --start --quiet --exec /usr/bin/xvfb-run -- /opt/zaproxy/zap.sh || true
    while [ $(nc -z -v localhost <%= node['zap']['port'] %>; echo $?) == "1" ]; do sleep 1; done 2>/dev/null
    echo "$NAME."
    ;;
  stop)
    echo -n "Stopping $DESC: "
    pkill -u <%= node['zap']['user'] %>
    echo "$NAME."
    ;;
  *)
    N=/etc/init.d/$NAME
    echo "Usage: $N {start|stop}" >&2
    exit 1
    ;;
esac

exit 0
